# supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock   ; (the path to the socket file)

[supervisord]
logfile=/tmp/supervisord.log       ; (main log file; default $CWD/supervisord.log)
logfile_maxbytes=50MB              ; (max main logfile bytes before rotation; default 50MB)
logfile_backups=10                 ; (num of main logfile rotation backups; default 10)
loglevel=info                      ; (log level; default info; others: debug,warn,trace)
pidfile=/tmp/supervisord.pid       ; (supervisord pidfile; default supervisord.pid)
nodaemon=false                     ; (start in foreground if true; default false)
minfds=1024                        ; (min. available startup file descriptors; default 1024)
minprocs=200                       ; (min. available process descriptors; default 200)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL for a unix socket

[program:my-node-app]
directory = /home/ec2-user/app  ; <-- Adjust to your application directory
command = npm run dev
autostart = true
autorestart = true
startsecs = 3
stdout_logfile = /var/log/my-node-app-stdout.log
stderr_logfile = /var/log/my-node-app-stderr.log
; user = ec2-user                       ; (uncomment if desired to run as ec2-user)
